# scripts/mcp-env/com.mcp.env.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcp.env.agent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/smian/dotfiles/scripts/set-mcp-env-system.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>StandardOutPath</key>
    <string>/tmp/mcp-env-agent.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/mcp-env-agent.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/smian/dotfiles</string>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
